# home/dot_config/private_dconf/autostart/flameshot-custom.desktop
[Desktop Entry]
Type=Application
Name=Flameshot wayland fixed
Exec=sh -c "QT_QPA_PLATFORM=wayland flameshot"
X-GNOME-Autostart-enabled=true
